# steam-backlog-enforcer.service
[Unit]
Description=Steam Backlog Enforcer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/steam-backlog-enforcer
ExecStart=/usr/bin/python3 -m steam_backlog_enforcer.main enforce
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/steam-backlog-enforcer
Environment=HOME=/home/kuhy
# Hardening: enforcer must not be easily killed.
OOMScoreAdjust=-900

[Install]
WantedBy=multi-user.target
